# configfiles/hc952/files/etc/uci-defaults/03-tzone
#!/bin/sh

uci set system.@system[0].hostname='Huastlink'
uci commit system
/etc/init.d/system restart